# templates/systemd/group.target
[Unit]
Description=DESCRIPTION

[Install]
WantedBy=multi-user.target
